# etc/supervisor/conf.d/tcollector.conf
[program:tcollector]
command=/opt/tcollector/tcollector.py --host 10.10.11.14
